# files/etc/supervisor/conf.d/varnish.conf
[program:varnish]
command=/usr/sbin/varnishd -T :6082 -S /etc/varnish/secret -F -f /etc/varnish/default.vcl