# files/etc/supervisor/conf.d/varnishncsa.conf
[program:varnishncsa]
command=/varnishncsa.sh